# bin/bbs2_worker.service
[Unit]
Description=BBS2 WorkerMan Service
After=network.target
Wants=network.target

[Service]
Type=forking
ExecStart=/usr/bin/php /var/www/html/bbs2-api/bin/start_worker.php start -d
ExecReload=/usr/bin/php /var/www/html/bbs2-api/bin/start_worker.php restart -d
ExecStop=/usr/bin/php /var/www/html/bbs2-api/bin/start_worker.php stop
PrivateTmp=true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
